# etc/systemd/coolerdash.service
# -----------------------------------------------------------------------------
# Created by: damachine (christkue79 at gmail dot com)
# Website: https://github.com/damachine/coolerdash
# -----------------------------------------------------------------------------
# Info:
# 	CoolerDash systemd daemon service file
# 	This file is part of CoolerDash, a C99 LCD daemon for monitoring and controlling
#	Project coding standards and packaging notes (see README for details)
# Details:
#   Simple systemd service for CoolerDash daemon.
#   Runs as dedicated user 'coolerdash' for security.
#   Requires coolercontrol services to be running.
#   Uses configuration file at '/etc/coolerdash/config.ini'.
#   PID file at '/tmp/coolerdash.pid' for process management.
#   Do not run as root. Use dedicated user for security.
#   See 'README.md' and 'AUR-README.md' for details.
# Example:
#   To enable and start the service, use:
#   'sudo systemctl enable coolerdash.service'
#   'sudo systemctl start coolerdash.service'
#   To check the status of the service, use:
#   'systemctl status coolerdash.service'
# -----------------------------------------------------------------------------

[Unit]
Description=CoolerDash daemon
Documentation=man:coolerdash(1)
After=coolercontrold.service
Requires=coolercontrold.service

[Service]
Type=simple
User=coolerdash
ExecStartPre=/usr/bin/sleep 10
ExecStart=/usr/bin/coolerdash
PIDFile=/tmp/coolerdash.pid
RuntimeDirectory=coolerdash
ReadWritePaths=/tmp /opt/coolerdash
ProtectSystem=full
Restart=always
RestartSec=3
TimeoutStopSec=3

[Install]
WantedBy=multi-user.target
